# roles/wsgiapp/templates/wsgiapp.socket.j2
[Unit]
Description={{ socket_description }}

[Socket]
ListenStream={{ socket_path }}

[Install]
WantedBy=sockets.target
